# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:clickhouse]
command=/usr/local/bin/clickhouse-server --config-file=/etc/clickhouse-server/config.xml
directory=/var/lib/clickhouse
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/clickhouse.log
stderr_logfile=/var/log/supervisor/clickhouse-err.log
startsecs=5

[program:nats]
command=/usr/bin/nats-server --js --sd /var/lib/nats -m 8222
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nats.log
stderr_logfile=/var/log/supervisor/nats-err.log
startsecs=3

[program:redis]
command=/usr/bin/redis-server --save 60 1 --dir /var/lib/redis --daemonize no
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-err.log
startsecs=3

[program:aggregator]
command=/usr/local/bin/flowpulse-aggregator
autostart=true
autorestart=true
priority=20
startretries=10
startsecs=5
stdout_logfile=/var/log/supervisor/aggregator.log
stderr_logfile=/var/log/supervisor/aggregator-err.log

[program:server]
command=/usr/local/bin/flowpulse-server
autostart=true
autorestart=true
priority=20
startretries=10
startsecs=5
stdout_logfile=/var/log/supervisor/server.log
stderr_logfile=/var/log/supervisor/server-err.log

[program:web]
command=node /opt/flowpulse-web/server.js
directory=/opt/flowpulse-web
autostart=true
autorestart=true
priority=30
startretries=5
startsecs=5
environment=PORT="3000",HOSTNAME="0.0.0.0",NODE_ENV="production"
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web-err.log
